# systemd-services/ollama-pull.timer
[Unit]
Description=trigger ollama-pull user service on timer

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
